# files/debian/profile.d/scripts-path.sh
#!/bin/bash
export PATH=$PATH:/opt/scripts