# files/debian/profile.d/xdebug.sh
export XDEBUG_MODE=debug,develop,profile
export XDEBUG_SESSION=1
export PHP_IDE_CONFIG="serverName=g7system.local"